# contrib/juniper/autoreload.cron
*/15 * * * * /bin/sh /var/tmp/autoreload.sh
